# units/mysql@.service.d/publish-3306.conf
[Service]
# Expose the MySQL port (3306)
Environment='DOCKER_PUBLISH_3306=3306:3306'
